# admin/supervisor_config/lms.conf
[program:lms]
command = /venvs/lms/bin/gunicorn --log-file=- --workers 16 lms:app -b localhost:8000
directory = /ce-demo-lms