# admin/supervisor_config/proxy.conf
[program:proxy]
command = /proxy/cloud_sql_proxy -instances=<project-id>:<region>:<sql_name>=tcp:3306
directory = /proxy