# admin/supervisor_config/transcode.conf
[program:vtc]
command = /venvs/transcode/bin/python vtc.py
directory = /ce-demo-lms
